# com.ghostcat.rocha-tracker.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ghostcat.rocha-tracker</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/YOUR_USER/polymarket_rocha_extradition/run.sh</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>ROCHA_PHONE</key>
        <string>+1XXXXXXXXXX</string>
    </dict>

    <!-- Daily at 09:00 local time. Adjust Hour/Minute below. -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/YOUR_USER/polymarket_rocha_extradition/logs/launchd.out.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/YOUR_USER/polymarket_rocha_extradition/logs/launchd.err.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
